# conf/etc/service/statsd/run
#!/bin/sh

exec node /opt/statsd/stats.js /opt/statsd/config/$STATSD_INTERFACE.js 2>&1
